# packages/net-misc/dhcp/files/systemd/dhcpd.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=DHCP Daemon
Requires=network.target
After=network.target

[Service]
EnvironmentFile=/etc/conf.d/dhcpd.conf
ExecStart=/sbin/dhcpd -cf ${DHCPD_CONF} -q -f $DHCPD_OPTS ${DHCPD_IFACE}

[Install]
WantedBy=multi-user.target
