# packages/www-servers/lighttpd/files/systemd/lighttpd.service
[Unit]
Description=lighttpd web server
After=network.target

[Service]
ExecStart=/usr/sbin/lighttpd -f /etc/lighttpd/lighttpd.conf -D

[Install]
WantedBy=http-daemon.target
