# packages/www-servers/lighttpd/files/systemd/lighttpd.service
# Copyright 2010, 2011 Johannes Nixdorf <mixi@user-helfen-usern.de>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=lighttpd web server
After=network.target

[Service]
ExecStart=/usr/sbin/lighttpd -f /etc/lighttpd/lighttpd.conf -D

[Install]
WantedBy=http-daemon.target
