# packages/net/ntp/files/ntpd.service
[Unit]
Description=Network Time Service
After=syslog.target

[Service]
ExecStart=/usr/sbin/ntpd -n -g -p /var/run/ntpd.pid

[Install]
WantedBy=multi-user.target
