# packages/sys-process/vixie-cron/files/vixie-cron.service
[Unit]
Description=Command Scheduler
After=syslog.target

[Service]
ExecStart=/usr/sbin/cron -n

[Install]
WantedBy=multi-user.target
