# packages/net-fs/nfs-utils/files/systemd/proc-fs-nfsd.mount
# Copyright 2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=/proc/fs/nfsd

[Mount]
What=nfsd
Where=/proc/fs/nfsd
Type=nfsd
Options=noauto,nodev,noexec,nosuid
