# packages/net-fs/nfs-utils/files/systemd/rpcbind.service
[Unit]
Description=RPC Service
Before=rpcbind.target

[Service]
ExecStart=/sbin/rpcbind -d
RestartSec=1

[Install]
WantedBy=multi-user.target
