# packages/net-fs/nfs-utils/files/systemd/rpcbind.service
[Unit]
Description=RPC Service
After=syslog.target

[Service]
ExecStart=/sbin/rpcbind
Type=forking

[Install]
WantedBy=multi-user.target
