# packages/net-fs/nfs-utils/files/systemd/rpcgssd.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=rpcsec_gss daemon
Requires=rpcbind.service
After=rpcbind.service

[Service]
ExecStart=/usr/sbin/rpc.gssd -f

[Install]
WantedBy=network.target
