# packages/net-fs/nfs-utils/files/systemd/rpcgssd.service
[Unit]
Description=rpcsec_gss daemon
Requires=rpcbind.service
After=rpcbind.service

[Service]
ExecStart=/usr/sbin/rpc.gssd
Type=forking
