# packages/net-fs/nfs-utils/files/systemd/rpcidmapd.service
[Unit]
Description=NFSv4 ID <-> Name Mapper
Requires=rpcbind.service
After=rpcbind.service

[Service]
ExecStart=/usr/sbin/rpc.idmapd -f
