# packages/net-fs/nfs-utils/files/systemd/rpcmountd.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=NFS mount daemon
Requires=proc-fs-nfsd.mount proc-fs-nfs.mount rpcbind.service
After=proc-fs-nfsd.mount proc-fs-nfs.mount rpcbind.service

[Service]
ExecStart=/usr/sbin/rpc.mountd -F

[Install]
WantedBy=network.target multi-user.target
