# packages/net-fs/nfs-utils/files/systemd/rpcstatd.service
[Unit]
Description=NSM service daemon
Requires=rpcbind.service
After=rpcbind.service

[Service]
ExecStart=/sbin/rpc.statd
Type=forking

[Install]
WantedBy=multi-user.target
