# packages/net-fs/nfs-utils/files/systemd/rpcstatd.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=NSM service daemon
Requires=rpcbind.service
Before=rpcbind.target
After=rpcbind.service

[Service]
ExecStart=/sbin/rpc.statd -F
ExecStop=/bin/umount -a -f -t nfs,nfs4

[Install]
WantedBy=multi-user.target
