# packages/dev-lang/php/files/systemd/php-fpm.service
[Unit]
Description=php-fpm

[Service]
ExecStart=/usr/sbin/php-fpm

[Install]
WantedBy=multi-user.target
