# packages/sys-apps/rsyslog/files/systemd/rsyslog.service
[Unit]
Description=System Logging Service

[Service]
ExecStart=/usr/sbin/rsyslogd -n -c4
Restart=restart-always

[Install]
WantedBy=multi-user.target
