# packages/net-directory/openldap/files/slapd-initd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-nds/openldap/files/slapd-initd,v 1.2 2008/10/14 10:29:44 robbat2 Exp $

depend() {
    need net
}

start() {
    ebegin "Starting ldap-server"

    mkdir -p /run/openldap
    chown ldap:ldap /run/openldap
    chmod 755 /run/openldap

    eval start-stop-daemon --start --pidfile /run/openldap/slapd.pid --exec /usr/libexec/openldap/slapd -- -u ldap -g ldap "${OPTS}"
    eend ${?}
}

stop() {
    ebegin "Stopping ldap-server"
    start-stop-daemon --stop --signal 2 --quiet --pidfile /run/openldap/slapd.pid
    eend ${?}
}

