# packages/sys-boot/plymouth/files/systemd/plymouth-init.service
[Unit]
Description=Plymouth Initialization
After=local-fs.target

[Service]
Type=oneshot
ExecStart=-/usr/bin/plymouth --sysinit

[Install]
WantedBy=multi-user.target
