# packages/sys-boot/plymouth/files/systemd/plymouth-quit.service
[Unit]
Description=Plymouth Stop
After=displaymanager.service

[Service]
Type=oneshot
ExecStart=-/usr/bin/plymouth --quit

[Install]
WantedBy=graphical.target
