# packages/sys-apps/hal/files/systemd/hald.service
[Unit]
Description=Hardware abstraction layer
Requires=dbus.target
After=dbus.target

[Service]
Type=dbus
BusName=org.freedesktop.Hal
ExecStart=/usr/sbin/hald --daemon=no --use-syslog

[Install]
WantedBy=graphical.target
