# packages/sys-process/cronie/files/systemd/cronie.service
[Unit]
Description=Command Scheduler
After=syslog.target

[Service]
EnvironmentFile=/etc/conf.d/cronie.conf
ExecStart=/usr/sbin/crond -n $CRONDARGS

[Install]
WantedBy=multi-user.target
