# packages/sys-fs/udev/files/80-load-modules.rules
# This will automatically load any modules available for any devices
# seen by udev.
DRIVER!="?*", ENV{MODALIAS}=="?*", RUN+="/sbin/modprobe $ENV{MODALIAS}"
