# packages/net-misc/oidentd/files/oidentd
#!/sbin/runscript
# Copyright 2009 Saleem Abdulrasool <compnerd@compnerd.org>
# vim: set sw=4 sts=4 et :

depend() {
    need net
}

start() {
    ebegin "Starting oidentd"
    start-stop-daemon --start --quiet --exec /usr/sbin/oidentd -- -u nobody -g nobody
    eend $?
}

stop() {
    ebegin "Stopping oidentd"
    start-stop-daemon --stop --quiet --exec /usr/sbin/oidentd
    eend $?
}

