# packages/sys-sound/alsa-utils/files/systemd/alsasound-restore.service
[Unit]
Description=Restore ALSA Mixer Settings

[Service]
ExecStart=/usr/sbin/alsactl -g restore

[Install]
Also=alsasound-store.service
WantedBy=multi-user.target
