# packages/sys-sound/alsa-utils/files/systemd/alsasound-store.service
[Unit]
Description=Save ALSA Mixer Settings

[Service]
Type=oneshot
ExecStart=/usr/sbin/alsactl -g store
RemainAfterExit=yes

[Install]
Also=alsasound-restore.service
WantedBy=shutdown.target
