# packages/sys-sound/alsa-utils/files/systemd/alsasound-store.service
[Unit]
Description=Save ALSA Mixer Settings

[Service]
ExecStart=/usr/sbin/alsactl -g store

[Install]
Also=alsasound-restore.service
WantedBy=multi-user.target
