# packages/net-directory/openldap/files/systemd/slapd.service
[Unit]
Description=OpenLDAP server daemon
After=syslog.target

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/slapd.conf
ExecStart=/usr/libexec/openldap/slapd -u ldap -g ldap $OPTS

[Install]
WantedBy=multi-user.target
