# packages/sys-apps/systemd/files/basic.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=Miscellaneous basic configurations
DefaultDependencies=no
Conflicts=shutdown.target
Before=shutdown.target emergency.service emergency.target udev.service systemd-logger.socket
RefuseManualStart=yes
RefuseManualStop=yes

[Service]
Type=oneshot
StandardInput=tty
RemainAfterExit=yes
ExecStart=/bin/dmesg -n 1

[Install]
WantedBy=basic.target
