# packages/sys-apps/systemd/files/bootmisc.service
[Unit]
Description=Miscellaneous early userspace activities
Requires=local-fs.target
After=local-fs.target
Before=syslog.target
RefuseManualStart=yes
RefuseManualStop=yes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/etc/init.d/bootmisc.bash

[Install]
WantedBy=multi-user.target
