# packages/sys-apps/systemd/files/fsck.service
[Unit]
Description=fsck
DefaultDependencies=no
Conflicts=shutdown.target
Before=shutdown.target emergency.service emergency.target
RefuseManualStart=yes
RefuseManualStop=yes
Requires=remount-rootfs.service
Before=basic.target fsck.target

[Service]
Type=oneshot
StandardInput=tty
ExecStart=/LIBDIR/systemd/fsck.bash

[Install]
WantedBy=basic.target
