# packages/sys-apps/systemd/files/sysinit.service
[Unit]
Description=System Initialization
DefaultDependencies=no
Conflicts=shutdown.target
Before=shutdown.target emergency.service emergency.target
RefuseManualStart=yes
RefuseManualStop=yes
Before=remount-rootfs.service

[Service]
Type=oneshot
StandardInput=tty
RemainAfterExit=yes
ExecStart=/LIBDIR/systemd/sysinit.bash

[Install]
WantedBy=sysinit.target
