# packages/sys-fs/fuse/files/fuse.initd
#!/sbin/runscript
# Copyright 2009 Mike Kelly
# Distributed under the terms of the GNU General Public License v2

start() {
    ebegin "Loading fuse module"
    modprobe fuse
    eend $?

    ebegin "Mounting fuse control filesystem"
    mount -t fusectl fusectl /sys/fs/fuse/connections
    eend $?
}

stop() {
    ebegin "Unmounting fuse control filesystem"
    umount /sys/fs/fuse/connections
    eend $?

    ebegin "Unloading fuse module"
    modprobe -r fuse
    eend $?
}

