# packages/dev-db/mysql/files/systemd/mysql.service
[Unit]
Description=MySQL server (single instance)
After=syslog.target

# The bl-1 init script contains a mass of possible configurations and options to
# run several MySQL server instances, define nice-ness and configure everything
# up to different defaults for databases on an SSD.
# In most usecases, a single server instance is enough. That's what this service
# does; it starts a single server instance with sensible defaults.
# If you need more, provide a patch.

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /var/run/mysqld
ExecStartPre=/bin/chown mysql:mysql /var/run/mysqld
ExecStart=/usr/sbin/mysqld --defaults-file=/etc/mysql/my.cnf --basedir=/usr --datadir=/var/lib/mysql --pid-file=/var/run/mysqld/mysqld.pid --socket=/var/run/mysqld/mysqld.sock --user=mysql --log-error=/var/log/mysql/mysql.log -W

[Install]
WantedBy=multi-user.target
