# packages/sys-apps/sysstat/files/systemd/sysstat.service
[Unit]
Description=sysstat

[Service]
ExecStart=/usr/lib/sa/sadc -F -L 60 -

[Install]
WantedBy=multi-user.target
