# packages/sys-apps/syslog-ng/files/systemd/syslog-ng.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=System Logger Daemon

[Service]
ExecStart=/usr/sbin/syslog-ng -F
Restart=always

[Install]
WantedBy=multi-user.target
