# packages/net-dns/bind/files/systemd/named.service
[Unit]
Description=Domain Name Server
After=network.target

[Service]
ExecStart=/usr/sbin/named -f -u named -n 1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
