# packages/net-dns/bind/files/systemd/named.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=Domain Name Server
After=network.target

[Service]
ExecStart=/usr/sbin/named -f -u named -n 1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
