# packages/app-admin/logrotate/files/systemd/logrotate.service
[Unit]
Description=logrotate

[Service]
Type=oneshot
TimeoutSec=0
ExecStart=/usr/sbin/logrotate /etc/logrotate.conf
