# packages/net-dns/ddclient/files/systemd/ddclient.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=Dynamic DNS updater client
Requires=network.target

[Service]
Type=forking
User=ddclient
Group=ddclient
ExecStart=/usr/sbin/ddclient -pid /run/ddclient/ddclient.pid -file /etc/ddclient/ddclient.conf -syslog
PIDFile=/run/ddclient/ddclient.pid

[Install]
WantedBy=multi-user.target
