# packages/sys-apps/lm_sensors/files/systemd/fancontrol.service
[Unit]
Description=lm_sensors's fancontrol
After=lm_sensor.service

[Service]
ExecStart=/usr/sbin/fancontrol

[Install]
WantedBy=multi-user.target
