# packages/net-misc/rsync/files/systemd/rsyncd.service
[Unit]
Description=rsync daemon
After=syslog.target

[Service]
ExecStart=/usr/bin/rsync --daemon --no-detach

[Install]
WantedBy=multi-user.target
