# packages/mail-mta/sendmail/files/systemd/sendmail.service
[Unit]
Description=sendmail server
After=syslog.target

[Service]
ExecStartPre=-/usr/bin/newaliases > /dev/null 2>&1
ExecStartPre=/usr/bin/rm -f /var/spool/mqueue/xf*
ExecStart=/usr/sbin/sendmail -bD -q30m -L sm-mta

[Install]
WantedBy=multi-user.target
