# packages/mail-mta/sendmail/files/systemd/sendmail-client.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=sendmail client
After=syslog.target

[Service]
Type=forking
ExecStart=/usr/sbin/sendmail -Ac -q30m -L sm-cm
PIDFile=/var/spool/clientmqueue/sm-client.pid

[Install]
WantedBy=multi-user.target
