# packages/net-misc/dhcpcd/files/systemd/dhcpcd.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=DHCP client (dhcpcd)
Wants=network.target
Before=network.target
After=local-fs.target

[Service]
ExecStart=/sbin/dhcpcd -q -B

[Install]
WantedBy=multi-user.target
