# packages/net-misc/dhcpcd/files/systemd/dhcpcd.service
[Unit]
Description=DHCP client (dhcpcd)
Before=network.target
After=local-fs.target

[Service]
ExecStart=/sbin/dhcpcd -q -B

[Install]
WantedBy=network.target
