# packages/app-admin/ulogd/files/systemd/ulogd.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=ulog daemon
After=network.target

[Service]
Environment=TMPDIR=
ExecStart=/usr/sbin/ulogd
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
