# packages/net/ntp/files/systemd/ntpd.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=Network Time Service
After=syslog.target

[Service]
ExecStart=/usr/sbin/ntpd -n -g -p /var/run/ntpd.pid

[Install]
WantedBy=multi-user.target
