# packages/net/ntp/files/systemd/ntpd.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=Network Time Service
After=syslog.target

[Service]
Type=forking
ExecStart=/usr/sbin/ntpd -u ntp:ntp -g -p /run/ntpd.pid
PrivateTmp=true

[Install]
WantedBy=multi-user.target
