# packages/net/ntp/files/systemd/ntpd.service
# Copyright 2010-2012 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=Network Time Service
After=network.target nss-lookup.target

[Service]
Type=forking
ExecStart=/usr/sbin/ntpd -u ntp:ntp -g
PrivateTmp=true
ControlGroup=cpu:/
Restart=always

[Install]
WantedBy=multi-user.target
