# packages/dev-scm/git/files/systemd/git-daemon@.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=git daemon Per-connection server

[Service]
User=${GIT_USER}
Group=${GIT_GROUP}
StandardInput=socket
StandardOutput=socket
EnvironmentFile=/etc/conf.d/git-daemon.conf
ExecStart=/usr/bin/git daemon --inetd $GIT_OPTS
