# packages/dev-scm/git/files/systemd/git-daemon.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=git daemon
After=syslog.target

[Service]
EnvironmentFile=/etc/conf.d/git-daemon.conf
ExecStart=/usr/bin/git daemon --pid-file=/var/run/git-daemon.pid --syslog --user=nobody --group=nobody $GIT_OPTS

[Install]
WantedBy=multi-user.target
