# packages/sys-apps/kexec-tools/files/systemd/kexec.service
[Unit]
Description=kexec
DefaultDependencies=no
Before=reboot.service

[Service]
Type=oneshot
EnvironmentFile=/etc/conf.d/kexec.conf
ExecStart=/usr/sbin/kexec -l ${KERNEL} --reuse-cmdline

[Install]
WantedBy=shutdown.target
