# packages/net-misc/openntpd/files/systemd/openntpd.service
[Unit]
Description=OpenNTP Daemon

[Service]
EnvironmentFile=/etc/conf.d/openntpd
ExecStart=/usr/sbin/ntpd -d $NTPD_OPTS
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
